# Git/usr/etc/profile.d/gawk.csh
version https://git-lfs.github.com/spec/v1
oid sha256:9ba2af9853121df6dd5edff9254a2946e39eab8dc6513348fea28f739ee96648
size 1107
