# Git/etc/profile.d/aliases.sh
version https://git-lfs.github.com/spec/v1
oid sha256:4678d4939642b9a5b13e2f700a6968550f8e936738879a625479df5429a60e53
size 604
