# Git/etc/profile.d/env.sh
version https://git-lfs.github.com/spec/v1
oid sha256:255c64eecf1dbd188ef3452e73b9e099f3bb0e56bc23ff1493e4cc57f724d52f
size 440
