# Git/etc/profile.d/git-prompt.sh
version https://git-lfs.github.com/spec/v1
oid sha256:f8afd96004c02ce0127da5d31ce47ef476e04acc9ea29e979397c0eb2cdb6502
size 1320
