# Git/etc/profile.d/lang.sh
version https://git-lfs.github.com/spec/v1
oid sha256:0a3a3b131ace34f11f428118dfe81b34da148e29b6bea3b027d79bebd47141a7
size 822
