# Git/etc/profile.d/perlbin.sh
version https://git-lfs.github.com/spec/v1
oid sha256:91bf8a45a3f04f51fe0eca146a71e278d41dbde19ba60dabc1e41c4639477c07
size 667
