# Git/etc/profile.d/perlbin.csh
version https://git-lfs.github.com/spec/v1
oid sha256:478a52dc440d5ca324f4c95238b39cc8c90d5d7d8d5601fa085ffef79fea04f4
size 699
